# systemd/odoo-config2git.timer
[Unit]
Description=Run task every 60 minutes

[Timer]
OnUnitActiveSec=60min
Persistent=true

[Install]
WantedBy=timers.target

